# geneweb.conffiles
/etc/init.d/geneweb
/etc/logrotate.d/geneweb
